# ansible/roles/minecraft-dedicated-server/files/minecraft.cron.j2
# Minecraft backup worlds every hour
0 * * * * {{ minecraft_bin_dir }}/mscs backup {{ minecraft_default_world }}
